# src/transformers_ocr.service
[Unit]
Description=Transformers OCR daemon
PartOf=graphical-session.target
After=graphical-session.target

[Service]
ExecStart=/usr/bin/transformers_ocr start --foreground
Restart=on-failure
RestartSec=5
Slice=session.slice

[Install]
WantedBy=graphical-session.target
